# home/crontab
* * * * * /bin/test -e /home/j/.config/check-mail && /usr/local/bin/mbsync -a
